# src/templates/pipedaladmind.service.template
[Unit]
Description=${DESCRIPTION}
After=network.target

[Service]
ExecStart=${COMMAND} 
Restart=always
RestartSec=60
WorkingDirectory=/var/pipedal
TimeoutStopSec=15
KillMode=process


[Install]
WantedBy=multi-user.target
